# examples/octf.service
[Unit]
Description=openclaw-team-in-feishu orchestrator daemon (OpenClaw agents in Feishu groups)
After=network-online.target
Wants=network-online.target
# If you run openclaw-gateway as a user systemd unit, this won't help directly;
# the daemon does its own preflight against `openclaw channels status --probe`
# at startup so a healthy gateway is required either way.

[Service]
Type=simple
User=root
WorkingDirectory=/etc/octf
# Secrets: put your env vars in a separate file readable only by root:
#   /etc/octf/secrets.env
#   PM_HOST_SECRET=...
#   DEV_SECRET=...
#   ...
# Then chmod 600 it. systemd loads the file before ExecStart.
EnvironmentFile=/etc/octf/secrets.env

# Daemon logs to stderr/stdout — systemd captures into journal.
# Use `octf daemon logs` (or journalctl -u octf -f) to tail.
# /usr/bin/octf is created by `npm link`. If your npm prefix is elsewhere,
# replace with the absolute path you get from `which octf`.
ExecStart=/usr/bin/octf daemon start --config /etc/octf/octf.json

# Restart on crash but back off to avoid tight loops on bad config.
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Hardening (optional — relax if you need broader fs access).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log /root/oc/.shared
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
